# .config/systemd/user/outlook-calendar-sync.service
[Unit]
Description=Outlook calendar synchronization service

[Service]
Type=oneshot
ExecStart=%h/scripts/outlook-calendar-sync.sh
