# .config/systemd/user/sync.gopass.service
[Unit]
Description=Sync gopass

[Service]
Type=oneshot
ExecStart=/usr/bin/gopass sync
