# .config/systemd/user/todo-sync.service
[Unit]
Description=Todo git repo synchronization service

[Service]
Type=oneshot
ExecStart=bash -c %h/scripts/todo-sync.sh
